# server-nfs/filesystem/etc/supervisor.d/supervisor-nfsd.ini
[program:nfsd]
command=nfsd
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=10
startretries=0
